# rockie_workspace/devel/etc/catkin/profile.d/05.catkin-test-results.sh
# generated from catkin/cmake/env-hooks/05.catkin-test-results.sh.develspace.in

export CATKIN_TEST_RESULTS_DIR="/home/bryant/wpi_sample_return_robot_challenge/rockie_workspace/build/test_results"
export ROS_TEST_RESULTS_DIR="$CATKIN_TEST_RESULTS_DIR"
